# rocon/devel/etc/catkin/profile.d/15.rocon_test.zsh
# generated from rocon_test/env-hooks/15.rocon_test.zsh.em

. "/home/turtlebot/ese_team_project/rocon/src/rocon_multimaster/rocon_test/shells/env.zsh"
