# rocon/devel/etc/catkin/profile.d/15.rocon_launch.bash
# generated from rocon_launch/env-hooks/15.rocon_launch.bash.em

. "/home/turtlebot/ese_team_project/rocon/src/rocon_tools/rocon_launch/shells/env.bash"
